# packages/matrix-rssbot/matrix_rssbot-0.1.4.tar.gz/matrix_rssbot-0.1.4/rssbot.service
[Unit]
Description=RSSbot
Requires=network.target

[Service]
Type=simple
User=rssbot
Group=rssbot
WorkingDirectory=/opt/rssbot
ExecStart=/opt/rssbot/.venv/bin/python3 -um rssbot
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target